# packages/mocaccino/enable_flathub/files/enable_flathub.service
[Unit]
Description=Enable flathub
After=network.target flatpak.target

[Service]
ExecStartPre=/bin/bash -c 'until [ -f /var/lib/flatpak/repo/config ]; do sleep 1; done'
ExecStart=/usr/local/bin/enable_flathub.sh
Restart=always
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
