# packages/mocaccino/enable_flathub/files/enable_flathub.service
[Unit]
Description=Enable flathub
After=network.target network-online.target

[Service]
ExecStart=/usr/local/bin/enable_flathub.sh
Type=oneshot
RemainAfterExit=true
TimeoutSec=300
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
